# contrib/systemd/sov.socket
[Socket]
ListenFIFO=%t/sov.sock
SocketMode=0600

[Install]
WantedBy=sockets.target